# systemd/eshu-profiler.service
[Unit]
Description=ESHU System Profiler
Documentation=https://github.com/eshu-installer/eshu
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/eshu profile --refresh
User=root
StandardOutput=journal
StandardError=journal

# Cache directory
CacheDirectory=eshu
CacheDirectoryMode=0755

[Install]
WantedBy=multi-user.target
